# deploy/systemd/bultrain-backup.timer
[Unit]
Description=Run the BulTrain backup once a day

[Timer]
# Daily at 03:00 UTC (~06:00 Sofia). The data changes rarely — admin edits and
# the occasional GTFS import — so a daily snapshot is the right recovery point.
OnCalendar=*-*-* 03:00:00
# If the box was off at the scheduled time, run as soon as it is back.
Persistent=true
# Avoid every timer on the box firing at the same instant.
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
